# deploy/windows/sortai-task.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Task Scheduler import file for sortai-dashboard.
  Import via:  schtasks /Create /XML deploy\windows\sortai-task.xml /TN "sortai-dashboard"
  Or:          Task Scheduler GUI -> Action -> Import Task...

  Before importing, replace every occurrence of SORTAI_ROOT with the absolute path
  to your sortAI project folder (e.g. C:\Users\you\Documents\git\sortAI).

  Uses python.exe with Hidden=true so Task Scheduler creates a hidden console —
  no window appears, stdout/stderr work normally, and "End Task" terminates the process directly.
  The 60-second logon delay gives mapped network drives time to reconnect.
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Starts the sortAI dashboard (watch + no-browser) at user logon.</Description>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
      <!-- Delay gives mapped network drives time to reconnect before sortai validates paths -->
      <Delay>PT60S</Delay>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <!-- InteractiveToken: runs in the logged-on user's session (needed to access mapped drives) -->
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>true</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <!-- python.exe with Hidden=true: Task Scheduler creates a hidden console — no window,
           stdout/stderr work normally, End Task kills the process directly.
           Replace SORTAI_ROOT with the absolute path to your sortAI project folder. -->
      <Command>SORTAI_ROOT</Command>
      <Arguments>-m sortai dashboard --watch --no-browser</Arguments>
      <WorkingDirectory>SORTAI_ROOT</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
